# docker/evermemos/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; --------------------------------------------------------------------------
; MongoDB
; --------------------------------------------------------------------------
[program:mongodb]
command=/usr/bin/mongod --dbpath /data/mongodb --bind_ip 127.0.0.1 --quiet
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --------------------------------------------------------------------------
; Elasticsearch
; --------------------------------------------------------------------------
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=elasticsearch
autostart=true
autorestart=true
environment=ES_JAVA_HOME="/usr/share/elasticsearch/jdk",ES_HOME="/usr/share/elasticsearch"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --------------------------------------------------------------------------
; Milvus standalone (embedded etcd, local storage)
; --------------------------------------------------------------------------
[program:milvus]
command=/opt/milvus/bin/milvus run standalone
environment=MILVUS_LOG_LEVEL="warn",ETCD_USE_EMBED="true",ETCD_DATA_DIR="/data/milvus/etcd",ETCD_CONFIG_PATH="/opt/milvus/configs/embedEtcd.yaml",COMMON_STORAGETYPE="local",LOCAL_STORAGEPATH="/data/milvus/storage"
directory=/opt/milvus
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --------------------------------------------------------------------------
; Redis
; --------------------------------------------------------------------------
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --dir /data/redis --appendonly yes --loglevel warning
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; --------------------------------------------------------------------------
; EverMemOS API server (waits for backing services)
; --------------------------------------------------------------------------
[program:evermemos]
command=/bin/bash -c 'sleep 15 && cd /app && uv run python src/run.py'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
